# deploy/site.service
[Unit]
Description=cameronbrooks.net Go site
Documentation=https://github.com/CameronBrooks11/cameronbrooks-site
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy
ExecStart=/home/deploy/site
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

# Runtime environment
Environment="SITE_ADDR=127.0.0.1:8080"
Environment="SITE_ENV=production"

# Optional: secrets drop-in (created manually on server after first deploy)
# File: /etc/systemd/system/site.service.d/env.conf
# Contents: [Service]\nEnvironment="SOME_SECRET=value"
EnvironmentFile=-/etc/systemd/system/site.service.d/env.conf

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
